# flexget/flexget.timer
[Unit]
Description=flexget

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=default.target
